# 2014-2015/RT2_91/Врублевский Александр/crontab
SHELL=/bin/bash

0-59/2 * * * *   /home/tols/cronmyscript.sh /home/tols/Music /home/tols/cronmyscriptlog
